# palace-daemon.service
[Unit]
Description=palace-daemon — MemPalace HTTP/MCP gateway
After=network.target

[Service]
Type=simple
User=radu
WorkingDirectory=/home/radu/palace-daemon
ExecStartPre=-/usr/bin/fuser -k 8085/tcp
ExecStartPre=-/usr/bin/rm -f /tmp/palace-daemon-8085.lock
ExecStart=/home/radu/.local/share/pipx/venvs/mempalace/bin/python main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
